# systemd/openanime.service
[Unit]
Description=openAnime backend bridge (IR router -> WebSocket + HTTP UI)
# Wait for the arduino-router socket owner and the network before starting.
After=network-online.target arduino-router.service
Wants=network-online.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/openAnime
# A system service's default PATH omits /usr/local/bin, where yt-dlp installs --
# without this the scraper's `yt-dlp` call fails and /play hangs/errors.
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
# mpv and wmctrl need the X session's display. The service runs as the same user
# that owns the Openbox session; the autostart grants it access (xhost), so
# pointing at :0 is enough -- no XAUTHORITY juggling.
Environment=DISPLAY=:0
ExecStart=/usr/bin/python3 /home/user/openAnime/backend/app.py
# Streaming sites need outbound DNS/HTTPS; nothing here listens before the
# socket exists, and router_reader() already retries if the socket is late.
Restart=on-failure
RestartSec=3
# Surface backend logs via journalctl -u openanime.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
